# deploy/systemd/nexus-dl.service
[Unit]
Description=NexusDL distributed download node
Documentation=https://github.com/nexus-dl/nexus-dl
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nexus
Group=nexus
WorkingDirectory=/var/lib/nexus-dl
ExecStart=/usr/local/bin/nexus-dl serve -c /etc/nexus-dl/config.yaml
ExecReload=/bin/kill -HUP $MAINPID

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/nexus-dl
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=true
RestrictRealtime=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
LockPersonality=true
MemoryDenyWriteExecute=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
TasksMax=4096
TimeoutStopSec=20s

# Restart policy
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60s
StartLimitBurst=3

# Environment
Environment=NEXUS_NODE_DATA_DIR=/var/lib/nexus-dl
EnvironmentFile=-/etc/nexus-dl/nexus-dl.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexus-dl

[Install]
WantedBy=multi-user.target
